# supervisor_services.conf
[program:unreal]
user=ircd
priority=50
autostart=true
autorestart=true
command=/home/ircd/unrealircd/bin/unrealircd -F

[group:irc]
programs=unreal
